# upstart/agentio-control.conf
#
# install with:
# sudo initctl start agentio-control
#
start on runlevel [2345]
# setuid control
chdir /home/control
exec ./controld -p 2010
respawn
respawn limit 10 90
